# so2-tp3-sistemas-ebebidos-serSalas/build/usuariostp3.service
[Unit]
Description=Servicio de Usuarios
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/users_tp3
User=skaliver

[Install]
WantedBy=multi-user.target
